# basecamp/files/basecamp_e5a3db7/basecamp/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:revocation]
command=/usr/local/bin/revocation
environment=PORT="8081",DATA_FILE="/var/lib/revocation/revoked.txt",GIN_MODE="release"
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:server]
command=/usr/local/bin/server
environment=PORT="80",DATABASE_URL="/var/lib/basecamp/basecamp.db",MIGRATIONS_DIR="/usr/local/share/basecamp/migrations",WEB_ROOT="/srv/web",REVOCATION_URL="http://127.0.0.1:8081",JWT_SECRET="%(ENV_JWT_SECRET)s",SLUG_SECRET="%(ENV_SLUG_SECRET)s",GIN_MODE="release"
autorestart=true
priority=20
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
